# scripts/runtime/com.vio83.mac-dev-optimizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vio83.mac-dev-optimizer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/padronavio/Projects/vio83-ai-orchestra/scripts/runtime/mac_dev_optimizer.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>14400</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/padronavio/Projects/vio83-ai-orchestra/data/logs/mac_optimizer_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/padronavio/Projects/vio83-ai-orchestra/data/logs/mac_optimizer_stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
