# automation/mac-scripts/com.vio83.real-max-maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vio83.real-max-maintenance</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/padronavio/Projects/vio83-ai-orchestra/scripts/runtime/run_real_max_maintenance.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>17</integer></dict>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>17</integer></dict>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/padronavio/Projects/vio83-ai-orchestra/automation/logs/real-max-maintenance-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/padronavio/Projects/vio83-ai-orchestra/automation/logs/real-max-maintenance-launchd-error.log</string>
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
